# templates/co.startupz.briefing.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>co.startupz.briefing</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>__PLUGIN_PATH__/scripts/generate-briefing.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>__HOUR__</integer>
		<key>Minute</key>
		<integer>__MINUTE__</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/startupz-briefing.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/startupz-briefing-error.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>__HOME__</string>
		<key>STARTUPZ_PLUGIN_PATH</key>
		<string>__PLUGIN_PATH__</string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
